# cb5/cb5-linaro-desktop-dp/overlay/etc/udev/rules.d/11-audio.rules
DEVPATH=="/devices/virtual/switch/h2w", RUN+="/bin/audio-judge.sh"


